# image-builder/crontab
* * * * * php /var/www/html/ltc/app/console app:fetch-exchange-rates > /tmp/cron.log